# claudecode/rootfs/etc/s6-overlay/s6-rc.d/init-claude/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Claude Code Add-on: Init oneshot
# ==============================================================================
exec /usr/bin/claude-init.sh
