# claudecode/rootfs/etc/s6-overlay/s6-rc.d/sshd/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Claude Code Add-on: sshd service
# ==============================================================================

if [ "${CLAUDE_SSH_ENABLED}" != "true" ]; then
    bashio::log.info "SSH disabled, sleeping..."
    exec sleep infinity
fi

# Check if the SSH port is actually mapped
if ! bashio::var.has_value "$(bashio::addon.port 22)"; then
    bashio::log.warning "SSH enabled but port 22 is not mapped — configure the port in add-on network settings"
    exec sleep infinity
fi

bashio::log.info "Starting SSH server on port 22..."
exec /usr/sbin/sshd -D -e
